# ops/ao-web.service.d/override.conf
[Service]
# Host-specific public URL for operator logs. Tailscale serve remains the
# security boundary and should continue forwarding 443 to 127.0.0.1:5173.
Environment=AO_WEB_PUBLIC_URL=https://mirrorborn.tailc1fd9.ts.net/
